# deploy/systemd/clawarcade-reviewer.service
[Unit]
Description=ClawArcade TopicLab reviewer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__DEPLOY_DIR__
EnvironmentFile=-__ENV_FILE__
ExecStart=/usr/bin/env python3 arcade_reviewer.py --loop --poll-interval 60 --repo-root __DEPLOY_DIR__ --registry-path __DEPLOY_DIR__/generated/reviewer_registry.json
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
